# kodi/upstart.conf
#
# WARNING:
# This file is under CM control - all manual changes will be removed
#

# kodi-upstart
# starts KODI on startup by using xinit.
# by default runs as kodi, to change edit below.
env USER=kodi

description     "KODI-barebones-upstart-script"
author          "Matt Filetto"

start on (filesystem and stopped udevtrigger)
stop on runlevel [016]

# tell upstart to respawn the process if abnormal exit
respawn
respawn limit 10 5
limit nice 21 21

script
  exec su -c "xinit /srv/bin/kodi -- -nolisten tcp -nocursor :0" $USER
end script
